# content/cron.d/sitecheck
# /etc/cron.d/sitecheck
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * root bash /nginx_sites.sh >> /var/log/cron.log